# bbappends/meta-nxp-zigbee-rcp/recipes-zigbee-rcp/recipes-zigbee-rcp-sdk/zigbee-rcp-sdk/imx8mm-jaguar-sentai/zb_config.service
[Unit]
Description=ZBOSS 019.240 Configuration Service
Documentation=man:zb_config.sh(8)
Before=zb_mux.service zb_app.service
After=network.target bluetooth.service
Wants=network.target
Requires=sysinit.target
ConditionPathExists=/usr/sbin/zb_config.sh

[Service]
Type=oneshot
ExecStart=/usr/sbin/zb_config.sh
RemainAfterExit=yes
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal
Environment="ZBOSS_VERSION=019.240"

# Enhanced error handling for ZBOSS 019.240
Restart=no
RestartSec=5
KillMode=process

# Security enhancements
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/local/zboss /tmp
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
Alias=zb_config.service
