# bbappends/meta-nxp-zigbee-rcp/recipes-zigbee-rcp/recipes-zigbee-rcp-sdk/zigbee-rcp-sdk/imx8mm-jaguar-sentai/zb_mux.service
[Unit]
Description=ZBOSS 019.240 Muxer Service
Documentation=man:zb_mux(8)
After=zb_config.service network.target
Wants=zb_config.service
Requires=sysinit.target
ConditionPathExists=/usr/sbin/zb_mux
ConditionPathExists=/usr/sbin/zb_mux.sh

[Service]
Type=forking
ExecStart=/usr/sbin/zb_mux.sh
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/var/run/zb_mux.pid
TimeoutStartSec=60
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
Environment="ZBOSS_VERSION=019.240"
EnvironmentFile=-/etc/default/zb_mux.env

# Enhanced restart policy for ZBOSS 019.240
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
LimitNOFILE=1024
LimitNPROC=512
MemoryMax=256M
CPUQuota=50%

# Security enhancements
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/local/zboss /tmp /var/run
PrivateTmp=yes
PrivateDevices=no
DeviceAllow=/dev/spidev3.0 rw
DeviceAllow=/dev/gpiochip1 rw
DeviceAllow=/dev/gpiochip3 rw

[Install]
WantedBy=multi-user.target
Alias=zb_mux.service
