# recipes-devtools/python/python3-improv/imx93-jaguar-eink/improv-eink.service
[Unit]
Description=Improv BLE WiFi onboarding service for imx93-jaguar-eink
Documentation=https://www.improv-wifi.com/ble/
# CRITICAL: Must start WITHOUT network - this service is used to CONFIGURE WiFi!
# Only require Bluetooth to be available, network is optional
After=bluetooth.target
Wants=bluetooth.target
# Do NOT require network.target - service must run even without network connection
# This allows users to configure WiFi when no network is available

[Service]
Type=simple
ExecStart=/usr/share/improv/onboarding-server-eink.py
Restart=always
RestartSec=12
# CRITICAL: Service must survive system sleep/wake cycles
# When PMIC is turned off (pm sleep --alloff), the entire system reboots on wake
# This service will automatically restart after reboot via Restart=always
# Board-specific environment variables for imx93-jaguar-eink
Environment="IMPROV_WIFI_INTERFACE=wlan0"
Environment="IMPROV_SERVICE_NAME=Improv-Eink"
Environment="IMPROV_CONNECTION_NAME=improv-eink"
# Optional: Override server host if needed
# Environment="IMPROV_SERVER_HOST=api.co.uk"

[Install]
WantedBy=default.target

