# recipes-bsp/eink-scheduler/eink-scheduler/eink-scheduler.service
[Unit]
Description=E-ink Display Scheduler
Documentation=https://github.com/DynamicDevices/eink-scheduler
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/eink-scheduler/eink-scheduler.conf

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/eink-scheduler -f /etc/eink-scheduler/eink-scheduler.conf
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eink-scheduler

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/eink-scheduler /tmp
PrivateTmp=true

# Resource limits
LimitNOFILE=1024
MemoryMax=64M

[Install]
WantedBy=multi-user.target

# Copyright (c) 2025 Dynamic Devices Ltd. All rights reserved.
# Maintainer: Alex J Lennon <ajlennon@dynamicdevices.co.uk>
# Support: info@dynamicdevices.co.uk
