# recipes-multimedia/pulseaudio/pulseaudio/imx8mm-jaguar-sentai/pulseaudio.service
[Unit]
Description=Sound Service

[Service]
User=root
ExecStart=/usr/bin/pulseaudio --system --log-target=journal
Restart=on-failure

[Install]
WantedBy=multi-user.target
